# buildbase/custom/common/etc/init.d/final3
#!/bin/sh
### BEGIN INIT INFO
# Provides:          final3
# Required-Start:    network
# Default-Start:     3
# Short-Description: out of the box level goal
### END INIT INFO
