# buildbase/custom/common/sbin/rctst/rc3
### BEGIN INIT INFO
# Provides:          rc3
# Required-Start:    rc5
# Restart-With:      rc4
# Default-Start:     3 5
# Default-Stop:      0 1 6
# Short-Description: fire up the network interfaces
### END INIT INFO
